# conf/asgard.conf
description "Netflix Asgard"

respawn
respawn limit 0 10

start on runlevel [2345]
stop on runlevel [!2345]

setuid ubuntu

chdir /home/ubuntu
exec /home/ubuntu/asgard/run
